# infra/hardware/systemd/edge-agent.service
[Unit]
Description=BlackRoad Edge Agent
After=docker.service

[Service]
Restart=always
ExecStart=/usr/bin/docker run --rm --gpus all -e GATEWAY_URL=http://127.0.0.1:30888 -e TS_TRUST=0.62 ghcr.io/blackroad/edge-agent:latest

[Install]
WantedBy=multi-user.target
