# infra/hardware/systemd/kiosk-browser.service
[Unit]
Description=BlackRoad Kiosk Browser
After=network-online.target

[Service]
Environment=URL=https://blackroad.network/lighthouse
ExecStart=/usr/bin/chromium --kiosk --noerrdialogs --incognito --disable-translate $URL
Restart=always

[Install]
WantedBy=graphical.target
